# cemc-rdp.conf
start on desktop-start
stop on desktop-end

script
	remmina -c cemc.remmina
end script
